# crontab
# tw2other
MAILTO="i@leaskh.com"
*/1 * * * * /usr/local/bin/php /home/leask/leaskh.com/t2os-leaskh/cron1.php
*/1 * * * * /usr/local/bin/php /home/leask/leaskh.com/t2os-syxnx/cron1.php
